# v1/squashfs-root/etc/inittab
::sysinit:/etc/rc.d/rcS
::respawn:/sbin/getty ttyS1 57600
::shutdown:/bin/umount -a
